# etc/profile.d/developerstudio.sh
DEVTOOLS=/opt/developerstudio12.6
PATH=$PATH:$DEVTOOLS/bin
export PATH DEVTOOLS
